# systemd/reticulumpi.service
[Unit]
Description=ReticulumPi Network Node
After=network-online.target rnsd.service
Wants=network-online.target rnsd.service

[Service]
Type=simple
User=reticulumpi
Group=reticulumpi
Environment=PATH=/opt/reticulumpi/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/reticulumpi/.venv/bin/reticulumpi --config /etc/reticulumpi/config.yaml
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/reticulumpi /home/reticulumpi/.reticulum /home/reticulumpi/.config/reticulumpi /home/reticulumpi/.local/share/reticulumpi /home/reticulumpi/.nomadnet /home/reticulumpi/.nomadnet-tui /opt/reticulumpi/meshchat/storage
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
